# airshell.service
[Unit]
Description=AirShell sensor daemon
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/home/airshell/airshell/venv/bin/python /home/airshell/airshell/daemon.py
WorkingDirectory=/home/airshell/airshell
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
